# root/etc/s6-overlay/s6-rc.d/init-config-amp-update/run
#!/usr/bin/with-contenv bash

# upgrade instances
OLD=$(cat /config/.ampdata/instances/Main/AMPVersion)

if [[ "${AMP_VERSION//./}" -gt "${OLD//./}" ]]; then
	echo "Upgrading all instances from ${OLD} to ${AMP_VERSION//./}"
	s6-setuidgid abc ampinstmgr UpgradeAll
fi

# cleanup after ourselves
rm -f /config/.ampdata/instances/AMPCache-*.zip